# watchman-web.service
[Unit]
Description=Watchman Web — Video browser
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=watchman
ExecStart=/usr/bin/python3 /opt/watchman/web.py --config /etc/watchman/watchman.conf
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
